# bbs.conf
[program:bbs]
command=/usr/local/bin/gunicorn wsgi --bind 0.0.0.0:2000 --workers 3 --worker-class gevent
directory=/home/ubuntu/bbs
autostart=true
autorestart=true
